# var/lib/dpkg/info/mysql-server-5.5.conffiles
/etc/logrotate.d/mysql-server
/etc/logcheck/ignore.d.paranoid/mysql-server-5_5
/etc/logcheck/ignore.d.server/mysql-server-5_5
/etc/logcheck/ignore.d.workstation/mysql-server-5_5
/etc/init.d/mysql
/etc/mysql/conf.d/mysqld_safe_syslog.cnf
/etc/mysql/debian-start
